# deploy/systemd/snapshot-gc.service
[Unit]
Description=TianGong LCA Snapshot Storage GC
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/projects/lca_workspace/tiangong-lca-calculator
EnvironmentFile=/home/ubuntu/projects/lca_workspace/tiangong-lca-calculator/.env
Environment=RUST_LOG=info
Environment=SNAPSHOT_GC_SNAPSHOT_RETENTION_DAYS=30
Environment=SNAPSHOT_GC_ORPHAN_RETENTION_DAYS=30
Environment=SNAPSHOT_GC_MAX_SNAPSHOTS=100
Environment=SNAPSHOT_GC_MAX_ORPHAN_DIRS=200
Environment=SNAPSHOT_GC_MAX_BYTES=2147483648
Environment=SNAPSHOT_GC_BATCH_SIZE=50
SyslogIdentifier=snapshot_gc
ExecStart=/home/ubuntu/projects/lca_workspace/tiangong-lca-calculator/target/release/snapshot_gc --execute
